# systemd/stock-tracker-scheduler.service
[Unit]
Description=Stock Tracker Scheduler (수급현황, AI추천)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/stock-tracker
ExecStart=/home/ubuntu/stock-tracker/venv/bin/python src/main.py --mode scheduler
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
